# systemd/.config/systemd/units/background.service
[Unit]
Description=Set background image

[Service]
Type=oneshot
ExecStart=%h/.config/i3/scripts/background next
