# systemd/.config/systemd/units/background.service
[Unit]
Description=Set background image

[Service]
Type=oneshot
ExecStart=/usr/bin/feh --no-fehbg --bg-scale --randomize %h/.config/backgrounds/
